# scripts/vmrootfs/etc/rc.local

# Send a signal to the host to stop the VMs
signal_vm_exit()
{
	nc -z 192.168.33.1 9882
}


